# test/99-Blackbriar.rules
SUBSYSTEM=="input", KERNELS=="event0", ACTION=="add", RUN+="/bin/systemctl start blackbriar-monitor.service"